# deploy/mixer.ini
# supervisord service definition for the Mixer app on Uberspace.
#
# Install once as ~/etc/services.d/mixer.ini (adjust paths/user if your
# Uberspace username differs from "mattmaxx"), then run:
#   supervisorctl reread
#   supervisorctl update
#   supervisorctl start mixer
#
# NODE_ENV=production makes Next.js load ~/mixer/.env.production automatically,
# so secrets do not need to be repeated here.

[program:mixer]
command=/home/mattmaxx/mixer/node_modules/.bin/next start -p 8723
directory=/home/mattmaxx/mixer
environment=NODE_ENV="production"
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/home/mattmaxx/logs/mixer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/mattmaxx/logs/mixer-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
